# init_scripts/generic/plone-standalone.service
[Unit]
Description=Plone
After=network.target
     
[Service]
Type=forking
ExecStart=/usr/local/Plone/zinstance/bin/plonectl start
ExecStop=/usr/local/Plone/zinstance/bin/plonectl stop
ExecReload=/usr/local/Plone/zinstance/bin/plonectl restart
   
[Install]
WantedBy=multi-user.target
